# launchd/com.ronen.bot-trad.tradebot-main.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ronen.bot-trad.tradebot-main</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-is</string>
    <string>__PROJECT_DIR__/.venv/bin/python</string>
    <string>-u</string>
    <string>__PROJECT_DIR__/tradeBot_main.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>ProcessType</key>
  <string>Background</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>PYTHONPATH</key>
    <string>__PROJECT_DIR__</string>
    <key>PATH</key>
    <string>__PROJECT_DIR__/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__PROJECT_DIR__/logs/launchd/tradebot-main.out.log</string>

  <key>StandardErrorPath</key>
  <string>__PROJECT_DIR__/logs/launchd/tradebot-main.err.log</string>
</dict>
</plist>
